# destiny_aria2/DEBIAN/postinst
#!/bin/sh
chmod 777 /etc/aria2/aria2.session
chown root:root -R /etc/aria2
chown root:root -R /etc/systemd/system/aria2.service
systemctl daemon-reload
systemctl enable aria2
