# clash-linux-amd64-1.18.0/DEBIAN/prerm
#!/usr/bin/env bash

service clash stop
systemctl disable clash

